# tm/ansible/cloudburstd.service
[Unit]
Description=cloudburstd

[Service]
ExecStart=/usr/local/bin/cloudburstd

[Install]
WantedBy=multi-user.target
